# .config/systemd/user/compton.service
[Unit]
Description=X compositor

[Service]
Type=simple
ExecStart=/usr/bin/compton

[Install]
WantedBy=default.target
