# .config/systemd/user/udiskie.service
[Unit]
Description=Udiskie disk management

[Service]
Type=simple
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
